# crontab
SHELL=/bin/bash
@reboot $HOME/iglvzx_ebooks/start.sh
