# crontab
SHELL=/bin/bash
@reboot $HOME/ebooks/start.sh
